# bash-scripting-suite/systemd/bash-scripts-backup.timer
[Unit]
Description=Daily backup timer for Bash Scripting Suite

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
